# system/overlays/config/etc/systemd/system/systemd-firstboot.service.d/override.conf
# Prompt everything
[Service]
ExecStart=systemd-firstboot --prompt
